# sec-camera.service
[Unit]
Description=Security Camera Dashboard
After=network.target

[Service]
Type=simple
User=pyr
WorkingDirectory=/home/pyr/sec-camera-dashboard
ExecStart=/home/pyr/sec-camera-dashboard/venv/bin/python -m uvicorn backend.main:app --host 0.0.0.0 --port 5001
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
